# systemd/system/sync_ssh_iam_users.timer
[Unit]
Description=sync_ssh_iam_users timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
