# systemd/system/sync_ssh_iam_users_update.timer
[Unit]
Description=sync_ssh_iam_users update timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
